# templates/etc/systemd/system/networking.service.d/override.conf.j2
# {{ ansible_managed }}
# ansibleguy.linux_networking

[Service]
ExecStartPost=/usr/bin/systemctl restart {{ startup_service }}
