# templates/etc/systemd/system/networking.service.d/override.conf.j2
[Service]
ExecStartPost=/usr/bin/systemctl restart {{ startup_service }}
